# rc.local
#!/bin/sh -e
#
# rc.local
#

cpufreq-set -r -u 580000
cpufreq-set -r -d 480000

echo 199 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio199/direction
echo 0 > /sys/class/gpio/gpio199/value

echo 203 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio203/direction
echo 0 > /sys/class/gpio/gpio203/value

echo 198 > /sys/class/gpio/export
echo out > /sys/class/gpio/gpio198/direction
echo 0 > /sys/class/gpio/gpio198/value

cd /srv/www
php -S 0.0.0.0:80 &

exit 0
